# .init/celeryd.conf
start on runlevel

respawn

script
    chdir /app/current/myapp
    exec /app/env/bin/python manage.py celeryd --settings settings
end script
