# .init/celeryd.conf
start on runlevel

chdir /app/current/myapp

exec /app/env/bin/python manage.py celeryd --settings settings
