# .init/celeryd.conf
start on runlevel

exec python /app/current/myapp/manage.py celeryd
